# firmware/boards/raspberry-pi-zero-2w/scripts/bootstrap/peripheral-init.service
[Unit]
Description=Peripheral Initialization
After=basic.target
Before=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/init-peripherals.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
